# ansible/roles/zfs_mirror/templates/zfs-monthly-scrub.service.j2
[Unit]
Description=ZFS monthly scrub

[Service]
Type=oneshot
ExecStart=/sbin/zpool scrub {{ zfs_pool_name }}
