# iac/ansible/roles/zfs_mirror/templates/zfs-monthly-scrub.timer.j2
[Unit]
Description=Run zfs-monthly-scrub timer

[Timer]
OnCalendar={{ zfs_scrub_oncalendar }}
Persistent=true

[Install]
WantedBy=timers.target
